# board/ifc6401/overlay/etc/systemd/system/fand.service
[Unit]
Description=Fand system fan management daemon
RequiresMountsFor=/persist

[Service]
ExecStart=/bin/fand

[Install]
WantedBy=multi-user.target
