# board/ifc6401/overlay/etc/udev/rules.d/30-ttys.rules
KERNEL=="ttyMSM*", RUN+="/bin/stty -F /dev/%k 115200"
